# resources/com.cwygoda.document-inbox.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cwygoda.document-inbox</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/cwygoda/.local/bin/uv</string>
        <string>run</string>
        <string>--project</string>
        <string>/Users/cwygoda/HobbyForge/document-inbox</string>
        <string>document-inbox</string>
        <string>watch</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/cwygoda/Library/Logs/document-inbox.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/cwygoda/Library/Logs/document-inbox.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/cwygoda/HobbyForge/document-inbox</string>
</dict>
</plist>
